# udev/00-usb-keyboards.rules
ACTION=="add", ATTRS{idVendor}=="feed", ATTRS{idProduct}=="6060", RUN+="/home/jaday/bin/kbd_udev", OWNER="jaday"
